# install/common/caddy/crontab
# Hourly regeneration of Caddyfile from DEEPWELL
*/5 * * * * /usr/local/bin/wikijump-update-caddyfile
