# install/common/caddy/crontab
# Hourly regeneration of Caddyfile from DEEPWELL
0 * * * * /usr/local/bin/wikijump-update-caddyfile
